# systemd/tv-downloader.timer
[Unit]
Description=Run TradingView Pine Script Downloader daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
